# systemd/file-sync.service
[Unit]
Description=GitHub File Sync Service
Documentation=https://github.com/rdone4425/github11
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/file-sync-system

[Service]
Type=forking
User=root
Group=root
WorkingDirectory=/file-sync-system
ExecStartPre=/file-sync-system/bin/file-sync validate
ExecStart=/file-sync-system/bin/file-sync-daemon start
ExecStop=/file-sync-system/bin/file-sync-daemon stop
ExecReload=/file-sync-system/bin/file-sync-daemon reload
PIDFile=/file-sync-system/logs/daemon.pid
Restart=always
RestartSec=10
RestartPreventExitStatus=2

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/file-sync-system/logs /file-sync-system/config

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 环境变量
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=LOG_LEVEL=INFO

[Install]
WantedBy=multi-user.target
